# csi/etc/services.d/wm/run
#!/bin/execlineb -P

/usr/bin/with-contenv /usr/bin/openbox
